# sisl-rag-api.service
[Unit]
Description=SISL RAG API
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/root/rag-api
Environment=PYTHONPATH=/home/root/rag-api
ExecStart=/home/root/rag-api/rag_env/bin/python api.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
